# node-agent/triangle-node.service
[Unit]
Description=Triangle Mesh node-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=triangle
WorkingDirectory=/opt/triangle/node-agent
Environment=TRIANGLE_CONFIG=/etc/triangle/config.yaml
ExecStart=/opt/triangle/node-agent/.venv/bin/python agent.py -c /etc/triangle/config.yaml
Restart=always
RestartSec=3
# Buffer survives restarts; store-and-forward means nothing is lost across a reboot.
StateDirectory=triangle

[Install]
WantedBy=multi-user.target
